# commons/DEBIAN/prerm
#!/bin/sh
set -e

# stop jasmind server
if [ -x "/etc/init.d/jasmind" ]; then
	invoke-rc.d jasmind stop || exit $?
fi
